# Django_blog.ini
[program:django_blog]
command=pipenv run gunicorn blogproject.wsgi -w 2 -k gthread -b 127.0.0.1:8000
directory=/home/yizhuo/apps/Django_blog
autostart=true
autorestart=unexpected
user=yizhuo
stdout_logfile=/home/yizhuo/etc/supervisor/var/log/django-blog-stdout.log
stderr_logfile=/home/yizhuo/etc/supervisor/var/log/django-blog-stderr.log
